# docker/supervisord.conf
# Enhanced Supervisord Configuration for LlamaAgent
# Author: Nik Jois <nikjois@llamasearch.ai>

[supervisord]
nodaemon=true
user=llamaagent
logfile=/app/logs/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/app/logs
loglevel=info
logfile_maxbytes=50MB
logfile_backups=10
silent=false
minfds=1024
minprocs=200

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700
chown=llamaagent:llamaagent

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

# Main Application Server
[program:llamaagent-api]
command=python -m uvicorn src.llamaagent.api.production_app:app --host 0.0.0.0 --port 8000 --workers 4 --worker-class uvicorn.workers.UvicornWorker --timeout-keep-alive 300
directory=/app
user=llamaagent
autostart=true
autorestart=true
startretries=3
stdout_logfile=/app/logs/api.log
stderr_logfile=/app/logs/api_error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_backups=10
environment=PYTHONPATH="/app/src",LLAMAAGENT_ENV="production"
priority=100
startsecs=10
stopwaitsecs=30
killasgroup=true
stopasgroup=true

# Background Worker
[program:llamaagent-worker]
command=python -m src.llamaagent.worker.main
directory=/app
user=llamaagent
autostart=true
autorestart=true
startretries=3
stdout_logfile=/app/logs/worker.log
stderr_logfile=/app/logs/worker_error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_backups=10
environment=PYTHONPATH="/app/src",LLAMAAGENT_ENV="production"
priority=200
startsecs=10
stopwaitsecs=30
killasgroup=true
stopasgroup=true

# Metrics Server
[program:llamaagent-metrics]
command=python -m src.llamaagent.monitoring.metrics_server --port 8001
directory=/app
user=llamaagent
autostart=true
autorestart=true
startretries=3
stdout_logfile=/app/logs/metrics.log
stderr_logfile=/app/logs/metrics_error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_backups=10
environment=PYTHONPATH="/app/src",LLAMAAGENT_ENV="production"
priority=300
startsecs=10
stopwaitsecs=30
killasgroup=true
stopasgroup=true

# Health Check Monitor
[program:llamaagent-health]
command=python -m src.llamaagent.monitoring.health_monitor
directory=/app
user=llamaagent
autostart=true
autorestart=true
startretries=3
stdout_logfile=/app/logs/health.log
stderr_logfile=/app/logs/health_error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_backups=10
environment=PYTHONPATH="/app/src",LLAMAAGENT_ENV="production"
priority=400
startsecs=10
stopwaitsecs=30
killasgroup=true
stopasgroup=true

# Log Rotation Service
[program:llamaagent-logrotate]
command=python -m src.llamaagent.utils.log_rotator
directory=/app
user=llamaagent
autostart=true
autorestart=true
startretries=3
stdout_logfile=/app/logs/logrotate.log
stderr_logfile=/app/logs/logrotate_error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_backups=10
environment=PYTHONPATH="/app/src",LLAMAAGENT_ENV="production"
priority=500
startsecs=10
stopwaitsecs=30
killasgroup=true
stopasgroup=true

# Database Maintenance
[program:llamaagent-db-maintenance]
command=python -m src.llamaagent.database.maintenance
directory=/app
user=llamaagent
autostart=true
autorestart=true
startretries=3
stdout_logfile=/app/logs/db_maintenance.log
stderr_logfile=/app/logs/db_maintenance_error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_backups=10
environment=PYTHONPATH="/app/src",LLAMAAGENT_ENV="production"
priority=600
startsecs=10
stopwaitsecs=30
killasgroup=true
stopasgroup=true

# Cache Cleaner
[program:llamaagent-cache-cleaner]
command=python -m src.llamaagent.cache.cleaner
directory=/app
user=llamaagent
autostart=true
autorestart=true
startretries=3
stdout_logfile=/app/logs/cache_cleaner.log
stderr_logfile=/app/logs/cache_cleaner_error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_backups=10
environment=PYTHONPATH="/app/src",LLAMAAGENT_ENV="production"
priority=700
startsecs=10
stopwaitsecs=30
killasgroup=true
stopasgroup=true

# Process group for core services
[group:llamaagent-core]
programs=llamaagent-api,llamaagent-worker,llamaagent-metrics,llamaagent-health
priority=100

# Process group for maintenance services
[group:llamaagent-maintenance]
programs=llamaagent-logrotate,llamaagent-db-maintenance,llamaagent-cache-cleaner
priority=200

# Event listeners
[eventlistener:llamaagent-crash-monitor]
command=python -m src.llamaagent.monitoring.crash_monitor
directory=/app
user=llamaagent
events=PROCESS_STATE_FATAL
autostart=true
autorestart=true
stdout_logfile=/app/logs/crash_monitor.log
stderr_logfile=/app/logs/crash_monitor_error.log
environment=PYTHONPATH="/app/src",LLAMAAGENT_ENV="production"

[eventlistener:llamaagent-memory-monitor]
command=python -m src.llamaagent.monitoring.memory_monitor
directory=/app
user=llamaagent
events=TICK_60
autostart=true
autorestart=true
stdout_logfile=/app/logs/memory_monitor.log
stderr_logfile=/app/logs/memory_monitor_error.log
environment=PYTHONPATH="/app/src",LLAMAAGENT_ENV="production"

# Include additional configuration files
[include]
files = /app/config/supervisor.d/*.conf 